# ansible/roles/upstart-conf/templates/celery.conf
# Managed by ansible, do not edit.
author "Milo Casagrande <milo.casagrande@linaro.org>"
description "Configuration to run celery at startup."

start on started redis-server
stop on stopping redis-server

respawn
respawn limit 10 5

kill timeout 6

limit nofile 4096 65536

setuid {{ apache_user }}
setgid {{ apache_user }}

chdir {{ install_base }}/{{ hostname }}/app

script
    exec celery worker --autoscale=4,1 --logfile=/var/log/celery/%h%I.log --loglevel=INFO --pidfile=/tmp/celery.pid --app=taskqueue
end script
