# ansible/roles/upstart-conf/templates/kernel-ci-backend.conf
# Managed by ansible, do not edit.
author "Milo Casagrande <milo.casagrande@linaro.org>"
description "Configuration to run tornado-base app at startup."

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5

limit nofile 4096 65536

setuid {{ app_user }}
setgid {{ app_user }}

chdir {{ install_base }}/{{ hostname }}/app

script
    if [ -d {{ install_base }}/.venv/{{ hostname }} ]; then
        . {{ install_base }}/.venv/{{ hostname }}/bin/activate
        {{ install_base }}/.venv/{{ hostname }}/bin/python -OO -R server.py
    else
        exec server.py
    fi
end script
