# ansible/roles/init-conf/templates/celery.conf
# Managed by ansible, do not edit.
author "Milo Casagrande <milo.casagrande@linaro.org>"
description "Configuration to run celery at startup."

start on started kernelci-backend
stop on stopping kernelci-backend

respawn
respawn limit 10 5

kill timeout 6

limit nofile 4096 65536

setuid {{ app_user }}
setgid {{ app_user }}

chdir {{ install_base }}/{{ hostname }}/app

script
    if [ -d {{ install_base }}/.venv/{{ hostname }} ]; then
        . {{ install_base }}/.venv/{{ hostname }}/bin/activate
        {{ install_base }}/.venv/{{ hostname }}/bin/celery worker --autoscale=4,1 --logfile=/var/log/celery/%h%I.log --loglevel=INFO --pidfile=/tmp/celery.pid --app=taskqueue
    else
        exec celery worker --autoscale=4,1 --logfile=/var/log/celery/%h%I.log --loglevel=INFO --pidfile=/tmp/celery.pid --app=taskqueue
    fi
end script
