# ansible/roles/init-conf/templates/kernelci-celery-beat.conf
# Managed by ansible, do not edit.
author "Milo Casagrande <milo.casagrande@linaro.org>"
description "Configuration to run celery at startup."

start on started kernelci-celery
stop on stopping kernelci-celery

respawn
respawn limit 10 5

kill timeout 6

limit nofile 4096 65536

setuid {{ app_user }}
setgid {{ app_user }}

chdir {{ install_base }}/{{ hostname }}/app

script
    if [ -d {{ install_base }}/.venv/{{ hostname }} ]; then
        {{ install_base }}/.venv/{{ hostname }}/bin/python -OO -R {{ install_base }}/.venv/{{ hostname }}/bin/celery beat --loglevel=INFO --schedule /var/run/celery/kernelci-beat.db --app=taskqueue
    else
        exec celery beat --loglevel=INFO --schedule /var/run/celery/kernelci-beat.db --app=taskqueue
    fi
end script
